# not-home/45-sda-power.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sda", ATTR{queue/rotational}=="1", RUN+="/usr/bin/hdparm -S 1 /dev/sda"
